# utils/freifunk-berlin-policy-routing/files/etc/uci-defaults/freifunk-berlin-policyroutting-ng
#!/bin/sh

# add routing tables
tables="/etc/iproute2/rt_tables"
grep -q "113 olsr-hnas" $tables || echo "113 olsr-hnas" >> $tables
